# tools/com.sakum.bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  launchd job: keeps the Sakum self-healing agent ALWAYS ON.

  It runs tools/sakum_bot.sh --once on a fixed interval (StartInterval=600s).
  Each pulse: fetches trusted PL release signals (read-only), generates a real
  sakum_lib_<topic>.s if missing, recompiles the whole x86-64 core, and on
  compile failure SELF-HEALS (rolls back bot-generated files, logs a mistake,
  exits 0 so launchd never crash-loops). On success it bumps the survivability
  counter in memory.md.

  The native socket server (tools/serve.s) is OPTIONAL and no longer required
  for the always-on behaviour; this plist drives the loop via the launchd
  timer directly. The bot is plain bash + native .s (doctrine-compliant).

  Install:
    cp tools/com.sakum.bot.plist ~/Library/LaunchAgents/
    launchctl load  ~/Library/LaunchAgents/com.sakum.bot.plist
  Stop:
    launchctl unload ~/Library/LaunchAgents/com.sakum.bot.plist
  Manual one-shot:
    bash tools/sakum_bot.sh --once
  -->
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.sakum.bot</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>/Users/Amit/Downloads/ Vendor Quotations/untitled folder 2/Sakum Lang/tools/sakum_bot.sh</string>
      <string>--once</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/Amit/Downloads/ Vendor Quotations/untitled folder 2/Sakum Lang</string>

    <key>StandardOutPath</key>
    <string>/Users/Amit/Downloads/ Vendor Quotations/untitled folder 2/Sakum Lang/tools/bot.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/Amit/Downloads/ Vendor Quotations/untitled folder 2/Sakum Lang/tools/bot.err.log</string>

    <!-- run the self-heal pulse on a fixed interval (the always-on trigger) -->
    <key>StartInterval</key>
    <integer>600</integer>

    <!-- do NOT keep-alive: a failing cycle must not crash-loop. Throttle
         restarts so launchd never spins. -->
    <key>KeepAlive</key>
    <false/>

    <key>ThrottleInterval</key>
    <integer>60</integer>

    <!-- start at login -->
    <key>RunAtLoad</key>
    <true/>
  </dict>
</plist>
