# tools/com.sakum.bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  launchd job: keeps the Sakum self-updater ALWAYS ON.
  It runs tools/serve.sh, which compiles tools/serve.s (raw x86-64 AT&T
  assembly) into native code and execs it. The server exposes the webhook +
  timer-pulse trigger endpoints; every pulse/webhook runs the bot
  (tools/sakum_bot.sh), which learns, generates real library functions,
  recompiles the core, and self-heals on failure.

  NOTE: there is deliberately NO Python host-language interpreter (serve.py is
  retired). SAKUM_LANG.md §2 requires every artifact to be raw assembly /
  machine code, not a foreign high-level runtime. tools/ are plain bash +
  native .s, which are permitted.

  Install:
    cp tools/com.sakum.bot.plist ~/Library/LaunchAgents/
    launchctl load  ~/Library/LaunchAgents/com.sakum.bot.plist
  Stop:
    launchctl unload ~/Library/LaunchAgents/com.sakum.bot.plist
-->
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.sakum.bot</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>/Users/Amit/Downloads/ Vendor Quotations/untitled folder 2/Sakum Lang/tools/serve.sh</string>
      <string>8080</string>
      <string>600</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/Amit/Downloads/ Vendor Quotations/untitled folder 2/Sakum Lang</string>

    <key>StandardOutPath</key>
    <string>/Users/Amit/Downloads/ Vendor Quotations/untitled folder 2/Sakum Lang/tools/bot.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/Amit/Downloads/ Vendor Quotations/untitled folder 2/Sakum Lang/tools/bot.err.log</string>

    <!-- keep-alive: relaunch if it ever exits -->
    <key>KeepAlive</key>
    <true/>

    <!-- start at login -->
    <key>RunAtLoad</key>
    <true/>
  </dict>
</plist>
